# wazuh/persistence/001_registry.xml
<!--
  HuntingThreats — Wazuh Hunt Rules
  Category:  Persistence — Registry (T1547, T1546)
  Batch:     1 / Persistence
  Rules:     101001 – 101025
  Platform:  Wazuh 4.x
  Sources:   Sysmon Event 12/13 (Registry key create/set)
             Sysmon Event 1 (Process Creation via reg.exe)
             Windows Event 4657 (Registry value modified — requires auditing)

  Persistence via Windows Registry covers:
    T1547.001 — Boot/Logon Autostart: Registry Run Keys
    T1546.010 — Event Triggered: AppInit DLLs
    T1546.011 — Event Triggered: Application Shimming
    T1547.004 — Winlogon Helper DLL
    T1547.005 — Security Support Provider
    T1547.012 — Print Processors

  Deployment note:
    Sysmon EventID 12/13: require Sysmon with registry monitoring configured.
    Field: win.eventdata.targetObject = full registry path.
    Field: win.eventdata.details = value set.
-->

<group name="hunt_persistence,registry,windows,">

  <!-- ── HKCU Run key — new value ─────────────────────────────────────────── -->

  <rule id="101001" level="10">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\</field>
    <description>Sysmon: HKCU Run key value set — autostart entry added for current user (T1547.001)</description>
    <mitre><id>T1547.001</id></mitre>
    <group>persistence,run_key,</group>
  </rule>

  <!-- ── HKLM Run key — new value ─────────────────────────────────────────── -->

  <rule id="101002" level="12">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\</field>
    <description>Sysmon: HKLM Run key value set — system-wide autostart entry added (T1547.001)</description>
    <mitre><id>T1547.001</id></mitre>
    <group>persistence,run_key,high_confidence,</group>
  </rule>

  <!-- ── RunOnce key ───────────────────────────────────────────────────────── -->

  <rule id="101003" level="10">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\CurrentVersion\\RunOnce(Ex)?\\</field>
    <description>Sysmon: RunOnce/RunOnceEx key set — one-time autostart on next logon (T1547.001)</description>
    <mitre><id>T1547.001</id></mitre>
    <group>persistence,run_key,</group>
  </rule>

  <!-- ── Run key value pointing to user-writable path ────────────────────── -->

  <rule id="101004" level="12">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\CurrentVersion\\Run(Once)?\\</field>
    <field name="win.eventdata.details" type="pcre2">(?i)(\\appdata\\|\\temp\\|\\downloads\\|\\desktop\\|\\users\\public\\)</field>
    <description>Sysmon: Run key pointing to user-writable path — high-risk autostart location (T1547.001)</description>
    <mitre><id>T1547.001</id></mitre>
    <group>persistence,run_key,high_confidence,</group>
  </rule>

  <!-- ── Image File Execution Options — debugger key ──────────────────────── -->

  <rule id="101005" level="14">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\Image File Execution Options\\.*\\Debugger</field>
    <description>Sysmon: IFEO Debugger set — execution hijacking via Image File Execution Options (T1546.012)</description>
    <mitre><id>T1546.012</id></mitre>
    <group>persistence,ifeo,execution_hijacking,critical,</group>
  </rule>

  <!-- ── Image File Execution Options — GlobalFlag silent exit ────────────── -->

  <rule id="101006" level="12">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\Image File Execution Options\\.*\\(GlobalFlag|MitigationOptions)</field>
    <description>Sysmon: IFEO GlobalFlag/MitigationOptions — process tampering for silent exit monitor (T1546.012)</description>
    <mitre><id>T1546.012</id></mitre>
    <group>persistence,ifeo,critical,</group>
  </rule>

  <!-- ── AppInit_DLLs ──────────────────────────────────────────────────────── -->

  <rule id="101007" level="14">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\Windows NT\\CurrentVersion\\Windows\\AppInit_DLLs</field>
    <description>Sysmon: AppInit_DLLs registry key set — DLL loaded into every GUI process (T1546.010)</description>
    <mitre><id>T1546.010</id></mitre>
    <group>persistence,appinitdlls,critical,</group>
  </rule>

  <!-- ── Winlogon Notify ───────────────────────────────────────────────────── -->

  <rule id="101008" level="14">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^12$|^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\Winlogon\\Notify\\</field>
    <description>Sysmon: Winlogon Notify key — DLL loaded at logon/logoff events (T1547.004)</description>
    <mitre><id>T1547.004</id></mitre>
    <group>persistence,winlogon,critical,</group>
  </rule>

  <!-- ── Winlogon Shell/Userinit replacement ───────────────────────────────── -->

  <rule id="101009" level="14">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\Winlogon\\(Shell|Userinit|Taskman)$</field>
    <description>Sysmon: Winlogon Shell/Userinit/Taskman modified — logon process hijacking (T1547.004)</description>
    <mitre><id>T1547.004</id></mitre>
    <group>persistence,winlogon,critical,</group>
  </rule>

  <!-- ── Security Support Provider (SSP) DLL ──────────────────────────────── -->

  <rule id="101010" level="14">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\LSA\\(Security Packages|Authentication Packages|Notification Packages)</field>
    <description>Sysmon: LSA Security/Auth/Notification Packages modified — SSP/SSPI DLL persistence (T1547.005)</description>
    <mitre><id>T1547.005</id></mitre>
    <group>persistence,lsa,ssp,critical,</group>
  </rule>

  <!-- ── Startup folder via registry ─────────────────────────────────────── -->

  <rule id="101011" level="12">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\User Shell Folders\\Startup|\\Shell Folders\\Startup</field>
    <description>Sysmon: Startup folder path modified via registry — redirected startup execution (T1547.001)</description>
    <mitre><id>T1547.001</id></mitre>
    <group>persistence,startup,high_confidence,</group>
  </rule>

  <!-- ── COM object hijacking (HKCU classes) ──────────────────────────────── -->

  <rule id="101012" level="12">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^12$|^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)HKCU\\SOFTWARE\\Classes\\CLSID\\</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)(InprocServer32|LocalServer32)</field>
    <description>Sysmon: HKCU CLSID InprocServer32 set — COM object hijacking for DLL loading (T1546.015)</description>
    <mitre><id>T1546.015</id></mitre>
    <group>persistence,com_hijacking,high_confidence,</group>
  </rule>

  <!-- ── BootExecute — custom value ───────────────────────────────────────── -->

  <rule id="101013" level="14">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\Session Manager\\BootExecute</field>
    <description>Sysmon: BootExecute registry key modified — pre-OS-load execution (T1547.001)</description>
    <mitre><id>T1547.001</id></mitre>
    <group>persistence,boot_execute,critical,</group>
  </rule>

  <!-- ── Browser Helper Object (BHO) ─────────────────────────────────────── -->

  <rule id="101014" level="10">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^12$|^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\Internet Explorer\\Extensions\\|\\Browser Helper Objects\\</field>
    <description>Sysmon: Browser Helper Object registry key — IE/browser extension persistence (T1176)</description>
    <mitre><id>T1176</id></mitre>
    <group>persistence,browser_extension,</group>
  </rule>

  <!-- ── Active Setup (per-user persistence) ──────────────────────────────── -->

  <rule id="101015" level="12">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^12$|^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\Active Setup\\Installed Components\\</field>
    <description>Sysmon: Active Setup component key created — per-user autostart on logon (T1547.014)</description>
    <mitre><id>T1547.014</id></mitre>
    <group>persistence,active_setup,high_confidence,</group>
  </rule>

  <!-- ── Scheduled Task via registry (run key created by schtasks) ─────────── -->

  <rule id="101016" level="10">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^12$|^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\Schedule\\TaskCache\\Tree\\</field>
    <description>Sysmon: Scheduled task registered in registry — persistence via task scheduler (T1053.005)</description>
    <mitre><id>T1053.005</id></mitre>
    <group>persistence,scheduled_task,</group>
  </rule>

  <!-- ── Print Processor DLL ──────────────────────────────────────────────── -->

  <rule id="101017" level="14">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^12$|^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\Control\\Print\\Environments\\.*\\Print Processors\\</field>
    <description>Sysmon: Print Processor DLL modified — DLL loaded by spooler on startup (T1547.012)</description>
    <mitre><id>T1547.012</id></mitre>
    <group>persistence,print_processor,critical,</group>
  </rule>

  <!-- ── Windows Terminal Profile persistence ─────────────────────────────── -->

  <rule id="101018" level="10">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\Microsoft\\Windows Terminal\\Profiles\\</field>
    <field name="win.eventdata.details" type="pcre2">(?i)(commandline|startingDirectory).*powershell|cmd|wscript</field>
    <description>Sysmon: Windows Terminal profile with shell command — possible persistence via terminal profile (T1547)</description>
    <mitre><id>T1547</id></mitre>
    <group>persistence,terminal_profile,</group>
  </rule>

  <!-- ── Netsh Firewall rule — inbound allow ──────────────────────────────── -->

  <rule id="101019" level="10">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^1$|^4688$</field>
    <field name="win.eventdata.image|win.eventdata.newProcessName" type="pcre2">(?i)\\netsh\.exe$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)firewall.*(add|allow|enable)|advfirewall.*add.*allow</field>
    <description>netsh.exe: firewall rule added — new inbound allow rule may expose service (T1562.004)</description>
    <mitre><id>T1562.004</id></mitre>
    <group>defense_evasion,firewall_tamper,</group>
  </rule>

  <!-- ── DisableAntiSpyware / DisableAntiVirus registry ───────────────────── -->

  <rule id="101020" level="14">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\Windows Defender\\(DisableAntiSpyware|DisableAntiVirus|DisableRealtimeMonitoring)</field>
    <description>Sysmon: Windows Defender disabled via registry — AV tampered to facilitate attack (T1562.001)</description>
    <mitre><id>T1562.001</id></mitre>
    <group>defense_evasion,av_tamper,critical,</group>
  </rule>

  <!-- ── Defender exclusion path added ────────────────────────────────────── -->

  <rule id="101021" level="12">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^12$|^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\Windows Defender\\Exclusions\\(Paths|Processes|Extensions)\\</field>
    <description>Sysmon: Windows Defender exclusion added — path/process excluded from scanning (T1562.001)</description>
    <mitre><id>T1562.001</id></mitre>
    <group>defense_evasion,av_tamper,high_confidence,</group>
  </rule>

  <!-- ── ETW (Event Tracing) provider disabled ────────────────────────────── -->

  <rule id="101022" level="12">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\HKLM\\SYSTEM\\CurrentControlSet\\Control\\WMI\\Autologger\\.*\\EnableLevel</field>
    <description>Sysmon: ETW autologger level changed — event tracing tampered to reduce visibility (T1562.006)</description>
    <mitre><id>T1562.006</id></mitre>
    <group>defense_evasion,etw_tamper,high_confidence,</group>
  </rule>

  <!-- ── SafeMode minimal services bypass ─────────────────────────────────── -->

  <rule id="101023" level="12">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^12$|^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\SYSTEM\\CurrentControlSet\\Control\\SafeBoot\\(Minimal|Network)\\</field>
    <description>Sysmon: SafeBoot registry modified — service added to run in safe mode (T1562.009)</description>
    <mitre><id>T1562.009</id></mitre>
    <group>defense_evasion,safeboot_tamper,high_confidence,</group>
  </rule>

  <!-- ── Terminal Server shadow key (RDP) ─────────────────────────────────── -->

  <rule id="101024" level="10">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\Terminal Server\\fSingleSessionPerUser|\\Terminal Services\\fDenyTSConnections</field>
    <description>Sysmon: RDP registry key modified — remote desktop settings changed (T1021.001)</description>
    <mitre><id>T1021.001</id></mitre>
    <group>lateral_movement,rdp_enable,</group>
  </rule>

  <!-- ── UserInitMprLogonScript ────────────────────────────────────────────── -->

  <rule id="101025" level="12">
    <if_group>windows</if_group>
    <field name="win.system.eventID" type="pcre2">^13$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)\\Environment\\UserInitMprLogonScript</field>
    <description>Sysmon: UserInitMprLogonScript set — script runs at every user logon (T1037.001)</description>
    <mitre><id>T1037.001</id></mitre>
    <group>persistence,logon_script,high_confidence,</group>
  </rule>

</group>
